# ha/monitor-cluster.cron
* * * * * root /usr/sbin/monitor-cluster
