# scripts/cleanup-bindmounts.init.in
#!/bin/bash
#
# cleanup-bindmounts: several services on this appliance might have bindmounts
#      which the halt script does not properly handle. If halt ever handles
#      bind mounts, this script can likely be removed.
#
# chkconfig: 2345 99 17
# description: Properly clean up bindmounts.
# config: @DESTDIR@/usr/bin

. /etc/init.d/functions

LOCKFILE=/var/lock/subsys/cleanup-bindmounts

case "$1" in
    start|restart)
        touch $LOCKFILE
    ;;
esac

if [ "$1" == "stop" ]
then
    for reps in 1 2; do for bindmount in `grep -E "[[:space:]]+bind[[:space:]]+" /etc/fstab | awk '{print $2}'`; do umount $bindmount 2>/dev/null; done; done
    rm -f $LOCKFILE
fi

