# scripts/job-server.init.in
#!/bin/bash
#
# rbuilder-isogen:  Starts the rBuilder ISO generation server.
#
# chkconfig: - 24 15
# description: This is the daemon which periodically requests a list of
#              ISO generation jobs from an rBuilder server and generates
#              the requested images.
# processname: job-server
# config: @DESTDIR@/usr/bin

# Source function library
. /etc/init.d/functions

# Configuration options
DEBUG_MODE=1    # set to non-zero to turn debug mode logging on

CONFIG_FILE='@DESTDIR@/srv/rbuilder/iso_gen.conf'
SCRIPT_PATH='@DESTDIR@/usr/share/rbuilder/scripts'
JOBSERVER_USER='isogen'
JOBSERVER_HOME='/srv/rbuilder/images'
JOBSERVER_SCRIPT="$SCRIPT_PATH/job-server"

LONGNAME="rBuilder ISO Generation Server"

#############################################################################

# Check for config file
if [ ! -r "$CONFIG_FILE" ]; then
    echo "Job server configuration file $CONFIG_FILE not found."
    exit 1
fi

# Check for debug mode and add to args if necessary
[ $DEBUG_MODE -gt 0 ] && EXTRA_ARGS="--debug $EXTRA_ARGS"

start() {
    echo -n "Starting $LONGNAME:"
    daemon --user $JOBSERVER_USER \
        HOME=$JOBSERVER_HOME \
        $JOBSERVER_SCRIPT --config $CONFIG_FILE $EXTRA_ARGS
    echo ""
}

stop() {
    echo -n "Stopping $LONGNAME:"
    killproc `basename $JOBSERVER_SCRIPT` -TERM
    echo ""
    if [ "$1" = "--wait" ]; then
        echo -n "Waiting for rBuilder ISO Generation server to terminate..."
        while [ -f "$LOCKFILE" ]; do sleep 1; done
        echo "done."
    fi 
}

status() {
    lockFile=`cat $CONFIG_FILE 2>&1 | grep "lockFile" | sed "s/.* //"`
    pid=`cat "$lockFile" 2>&1`
    case "$pid" in
        *'No such file or directory')
            echo "$LONGNAME is not running."
            ;;
        *'SHUTTING_DOWN')
            echo "$LONGNAME is shutting down."
            ;;
        *)
            case "`ps -p $pid -o comm=`" in
               'job-server')
                   echo "$LONGNAME is running."
                   ;;
               *)
                   echo "Lockfile is stale. Run '$0 stop' to remove it."
                   ;;
            esac
            ;;
    esac
}

case "$1" in
  start)
        start
	;;
  stop)
        stop
	;;
  restart)
        stop --wait
        start
	;;
  status)
        status job-server
	;;
  *)
    echo "Usage: `basename $0` {start|stop|restart|status}"
    ;;
esac
