# scripts/job-server.init.in
#!/bin/bash
#
# rbuilder-isogen:  Starts the rBuilder ISO generation server.
#
# chkconfig: - 24 15
# description: This is the daemon which periodically requests a list of
#              ISO generation jobs from an rBuilder server and generates
#              the requested images.
# processname: job-server 
# config: @DESTDIR@/usr/bin

CONFIG_FILE='@DESTDIR@/srv/mint/iso_gen.conf'
SCRIPT_PATH='@DESTDIR@/usr/share/mint/scripts'

case "$1" in
  start)
	$SCRIPT_PATH/job-server --config $CONFIG_FILE
	;;
  stop)
	$SCRIPT_PATH/job-server -k
	;;
  restart)
	$0 stop
	$0 start
	;;
  status)
	lockFile=`cat $CONFIG_FILE 2>&1 | grep "lockFile" | sed "s/.* //"`
	pid=`cat "$lockFile" 2>&1`
	case "$pid" in
	  *'No such file or directory')
		echo "Job server not running"
		exit 0
		;;
	  *)
		case "`ps -p $pid -o pid=`" in
		  '')
		    echo "lockfile is stale. fixing"
		    rm "$lockFile"
		    ;;
		  *)
		    echo "Job server is running"
		    ;;
		esac
		;;
	esac
	;;
  *)
    echo "Usage: `basename $0` {start|stop|restart|status}"
    ;;
esac
