# scripts/job-server.init.in
#!/bin/bash
#
# rbuilder-isogen:  Starts the rBuilder ISO generation server.
#
# chkconfig: - 24 15
# description: This is the daemon which periodically requests a list of
#              ISO generation jobs from an rBuilder server and generates
#              the requested images.
# processname: job-server
# config: @DESTDIR@/usr/bin

CONFIG_FILE='@DESTDIR@/srv/mint/iso_gen.conf'
SCRIPT_PATH='@DESTDIR@/usr/share/mint/scripts'
JOBSERVER_USER='isogen'

case "$1" in
  start)
        case `ps -A | grep "job-server"` in
        '')
            HOME=/home/$JOBSERVER_USER/ sudo -u $JOBSERVER_USER $SCRIPT_PATH/job-server --config $CONFIG_FILE
            ;;
        *) echo "Other job servers appear to be running. please stop them before continuing"
           ;;
        esac
	;;
  stop)
	HOME=/home/$JOBSERVER_USER/ sudo -u $JOBSERVER_USER $SCRIPT_PATH/job-server --config $CONFIG_FILE -k
	;;
  restart)
        $0 stop
        # Wait for job server to exit.
        lockFile=`cat $CONFIG_FILE 2>&1 | grep "lockFile" | sed "s/.* //"`
        echo "Waiting for rBuilder ISO Generation server to terminate."
        while [ -f "$lockFile" ]; do sleep 1; done
        $0 start
	;;
  status)
	lockFile=`cat $CONFIG_FILE 2>&1 | grep "lockFile" | sed "s/.* //"`
	pid=`cat "$lockFile" 2>&1`
	case "$pid" in
	  *'No such file or directory')
		echo "Job server is not running."
		;;
          *'SHUTTING_DOWN')
                echo "Job server is shutting down."
                ;;
	  *)
		case "`ps -p $pid -o comm=`" in
		  'job-server')
		    echo "Job server is running."
		    ;;
		  *)
		    echo "lockfile is stale. run '$0 stop' to remove it."
		    ;;
		esac
		;;
	esac
	;;
  *)
    echo "Usage: `basename $0` {start|stop|restart|status}"
    ;;
esac
