# scripts/multi-jobserver.init.in
#!/bin/bash
#
# rbuilder-isogen:  Starts multiple rBuilder ISO generation servers.
#
# chkconfig: - 24 15
# description: This executes multiple job server daemons in separate chroots.
# processname: job-server
# config: @DESTDIR@/usr/bin

JOBSERVER_BASEPATH='@DESTDIR@/%(servicedir)s/mint/jobserver'

case "$1" in
  start|stop|restart)
	for js in `ls $JOBSERVER_BASEPATH`
	  do
	  if [ -d "$JOBSERVER_BASEPATH/$js" ]
	      then
	      if [ "$1" == "start" ]
		  then
		  cp /etc/hosts "$JOBSERVER_BASEPATH/$js/etc/hosts"
		  cp /etc/resolv.conf "$JOBSERVER_BASEPATH/$js/etc/resolv.conf"
		  cp /etc/nsswitch.conf "$JOBSERVER_BASEPATH/$js/etc/nsswitch.conf"
	      fi
	      /usr/sbin/chroot "$JOBSERVER_BASEPATH/$js" /etc/init.d/rbuilder-isogen $1
	  fi
	done
	;;
    *)
	echo "Usage: `basename $0` {start|stop|restart}"
	;;
esac
