# distro/postgres.cron
0 0 * * * postgres /usr/bin/vacuumdb -azq -U rbuilder >& /dev/null

