# distro/postgres.cron
0 0 * * * postgres /usr/bin/vacuumdb -az -U rbuilder >& /dev/null

