# distro/rbuilder.cron
#m      h   DoM mon DoW     user    cmd
*/12    *   *   *   *       apache  /usr/share/rbuilder/scripts/update-package-index >& /dev/null
*/30    *   *   *   *       apache  /usr/share/rbuilder/scripts/update-package-index-external --cron >& /dev/null

*/5     *   *   *   *       root    /usr/share/rbuilder/scripts/job-cleanup -q
*       *   *   *   *       root    /usr/share/rbuilder/scripts/process-system-events -q >& /dev/null
