# distro/rbuilder.cron
#m      h   DoM mon DoW     user    cmd
*/12    *   *   *   *       apache  /usr/share/rbuilder/scripts/update-package-index >& /dev/null
*/30    *   *   *   *       apache  /usr/share/rbuilder/scripts/update-package-index-external --cron >& /dev/null

*/5     *   *   *   *       root    /usr/share/rbuilder/scripts/job-cleanup -q
# Capsule indexer
#*/10    *   *   *   *       root    curl -d "" --output /dev/null --stderr /dev/null http://127.0.0.1/api/capsules/rpm/content
*       *   *   *   *       root    /usr/share/rbuilder/scripts/process-system-events -q >& /dev/null
0       0   *   *   *       root    /usr/share/rbuilder/scripts/target-systems-import -q >& /dev/null
0       0   *   *   *       root    /usr/share/rbuilder/scripts/purge-surveys -q >& /dev/null
