# distro/rbuilder.cron
#m      h   DoM mon DoW     user    cmd
*/12    *   *   *   *       apache  /usr/share/rbuilder/scripts/update-package-index >& /dev/null
*/30    *   *   *   *       apache  /usr/share/rbuilder/scripts/update-package-index-external --cron >& /dev/null

28      3   *   *   *       root    /usr/share/rbuilder/scripts/rss-update -q
